# m2_demo/bash/profile.d/00-paths.sh
# Extend the PATH variable here.
#
# PATH="/opt/cool-software/bin/spin:$PATH"
# PATH="/opt/cool-script/bin/run:$PATH"
#
# export PATH
#

