# profile.d/40-man.sh
export MANPATH="$MANPATH:$__N2_DIR/man"
